# docker/crontab
# Health Hub scheduled jobs — runs inside hhub-cron container.
# Cron daemon must run as root; jobs drop to the hhub user.
# Times are container-local (set TZ in .env to match your locale, e.g. Europe/Moscow).

SHELL=/bin/sh
PATH=/opt/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# 21:00 — daily report (skips Fitbit collection if FITBIT_CLIENT_ID is unset).
# Reads whatever the Health Connect Bridge app pushed today + CPAP/O2Ring,
# formats a Telegram digest. With Fitbit disabled this is mostly a no-op
# until the formatter learns to surface hc_records — leave it on so cron
# logs prove the path stays alive.
0 21 * * * root runuser -u hhub -- /opt/venv/bin/hhub daily >> /app/logs/daily.log 2>&1

# 03:00 — online SQLite snapshot (WAL-safe), must not overlap the 21:00 window
0  3 * * * root runuser -u hhub -- /opt/venv/bin/hhub backup >> /app/logs/backup.log 2>&1

# --- Disabled (Fitbit Web API deprecation in Sept 2026; new app registrations closed) ---
# 0 10 * * * root runuser -u hhub -- /opt/venv/bin/hhub backfill --source cpap >> /app/logs/cpap.log 2>&1
#
# Note: `hhub backfill` is Fitbit-only in the current code (the --source flag
# is a stub). Re-enable when CPAP/O2Ring backfill is wired into the command.

# Required trailing newline for /etc/cron.d entries — keep this line.
